# etc/supervisord.d/neo4j.ini
[program:neo4j]
command=/opt/qnib/neo4j/bin/start.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart = false
autostart = true
startsecs=5
exitcodes=0
